# logrotate.cron.daily
#!/bin/sh

/usr/sbin/logrotate -d /etc/logrotate.conf > /var/log/logrotate.log
